# crontab
## crontab
#
MAILTO=admin@example.org

# renew the certificate
42 0 * * * zimbra /opt/letsencrypt-zimbra/obtain-and-deploy-letsencrypt-cert.sh -q
